# backbone/nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
user=root

[program:php5-fpm]
command=/usr/sbin/php5-fpm -F
user=root
autostart=true
